# supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -k start

[program:my-program-name]
command = /opt/start.sh
startsecs = 0
autorestart = false
startretries = 1

